# etc/supervisor/module.conf.in
[program:@MODULE_TARGET_NAME@]
autostart=false
autorestart=true
directory=@MOHSES_RUNTIME_DIR@
command=@CMAKE_INSTALL_PREFIX@/@MODULE_TARGET_NAME@ @MODULE_ARGUMENTS@
priority=105
numprocs=1
user=@MOHSES_USERNAME@
redirect_stderr=true
stdout_logfile=@MOHSES_LOGGING_DIR@/@MODULE_TARGET_NAME@.log
environment=LD_LIBRARY_PATH=@CMAKE_INSTALL_PREFIX@/bin:@CMAKE_INSTALL_PREFIX@/lib
