# etc/supervisor/module.conf.in
[program:@MODULE_TARGET_NAME@]
autostart=false
autorestart=true
directory=@MOHSES_RUNTIME_DIR@
command=@CMAKE_INSTALL_PREFIX@/@MODULE_TARGET_NAME@ @MODULE_ARGUMENTS@
priority=105
numprocs=1
user=@MOHSES_USERNAME@
redirect_stderr=true
stdout_logfile=@MOHSES_LOGGING_DIR@/amm_module_manager.log

